# mail/.config/systemd/user.control/my-get-mail.service
[Unit]
Description=my-get-mail

[Service]
Type=simple
ExecStart=%h/.local/bin/my-get-mail-daemon
Restart=always

[Install]
WantedBy=default.target
