# i3/.config/systemd/user.control/i3-my-kdeconnect-monitor.service
[Unit]
Description=i3-my-kdeconnect-monitor daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/i3-my-kdeconnect-monitor
Restart=always

[Install]
WantedBy=default.target
